# data/udisks-automount-agent.service.in
[Unit]
Description=Udisks user session automount/umount agent

[Service]
ExecStart=@udisksagentbindir@/udisks-automount-agent
Type=simple
Restart=always
RestartSec=1
Requires=udisks2.service
Requires=org.freedesktop.PolicyKit1.service
