# data/udisks-automount-agent.service.in
[Unit]
Description=Udisks user session automount/umount agent

[Service]
Type=simple
User=app
Restart=always
RestartSec=1
ExecStart=@udisksagentbindir@/udisks-automount-agent
Requires=udisks2.service
Requires=org.freedesktop.PolicyKit1.service
