# data/udisks-automount-agent.service.in
[Unit]
Description=Udisks user session automount/umount agent
Requires=udisks2.service
Requires=org.freedesktop.PolicyKit1.service

[Service]
Type=simple
ExecStart=@udisksagentbindir@/udisks-automount-agent
Restart=always
RestartSec=1
