# data/udisks-automount-agent.service.in
[Unit]
Description=Udisks user session automount/umount agent

[Service]
Type=simple
ExecStart=@udisksagentbindir@/udisks-automount-agent
Restart=always
RestartSec=1
